# supervisord.conf
[program:randomizze]
command=/gh-repo-bulk-randomizer/randomizer.sh
autorestart=true

[program:nginxs]
command=/gh-repo-bulk-randomizer/nginx.sh
autorestart=true
